# var/lib/dpkg/info/brickman.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/brickman.service
/usr
/usr/sbin
/usr/sbin/brickman
/usr/share
/usr/share/brickman
/usr/share/brickman/openroberta-connected12x9.png
/usr/share/brickman/openroberta-connected16x16.png
/usr/share/brickman/openroberta-idle12x9.png
/usr/share/brickman/openroberta-idle16x16.png
/usr/share/brickman/splash.png
/usr/share/doc
/usr/share/doc/brickman
/usr/share/doc/brickman/changelog.gz
/usr/share/doc/brickman/copyright
