# var/lib/dpkg/info/ev3-systemd.list
/.
/lib
/lib/modules-load.d
/lib/modules-load.d/ev3.conf
/lib/systemd
/lib/systemd/system
/lib/systemd/system/ev3-bluetooth.service
/lib/systemd/system/ev3-leds.service
/lib/systemd/system/ev3-uart@.service
/lib/systemd/system/ev3-usb@.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/60-ev3.rules
/usr
/usr/share
/usr/share/doc
/usr/share/doc/ev3-systemd
/usr/share/doc/ev3-systemd/changelog.gz
/usr/share/doc/ev3-systemd/copyright
/usr/share/ev3-systemd
/usr/share/ev3-systemd/ev3-usb.sh
