# var/lib/dpkg/info/openrobertalab.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/openrobertalab.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/openrobertalab.service
/usr
/usr/bin
/usr/bin/openrobertalab
/usr/lib
/usr/lib/python3
/usr/lib/python3/dist-packages
/usr/lib/python3/dist-packages/roberta
/usr/lib/python3/dist-packages/roberta/ter-u14n_unicode.pil
/usr/lib/python3/dist-packages/roberta/ter-u14n_unicode.pbm
/usr/lib/python3/dist-packages/roberta/ter-u18n_unicode.pbm
/usr/lib/python3/dist-packages/roberta/__init__.py
/usr/lib/python3/dist-packages/roberta/StaticData.py
/usr/lib/python3/dist-packages/roberta/ter-u12n_unicode.pbm
/usr/lib/python3/dist-packages/roberta/__version__.py
/usr/lib/python3/dist-packages/roberta/ev3.py
/usr/lib/python3/dist-packages/roberta/test_ev3.py
/usr/lib/python3/dist-packages/roberta/test.py
/usr/lib/python3/dist-packages/roberta/ter-u12n_unicode.pil
/usr/lib/python3/dist-packages/roberta/test_lab.py
/usr/lib/python3/dist-packages/roberta/lab.py
/usr/lib/python3/dist-packages/roberta/ter-u18n_unicode.pil
/usr/lib/python3/dist-packages/roberta/test_BlocklyMethods.py
/usr/lib/python3/dist-packages/roberta/BlocklyMethods.py
/usr/lib/python3/dist-packages/openrobertalab-1.7.4.egg-info
/usr/share
/usr/share/doc
/usr/share/doc/openrobertalab
/usr/share/doc/openrobertalab/changelog.gz
/usr/share/doc/openrobertalab/copyright
